# opt/arcOS/configs/udev/99-boot-device-watch.rules
ACTION=="remove", ENV{ID_FS_UUID}=="XXXARCHIVEUUIDXXX", RUN+="/opt/arcOS/bin/emergency-shutdown"
